# packaging/deb/etc/systemd/system/xhelix.service
[Unit]
Description=xhelix runtime security agent
Documentation=https://github.com/xhelix/xhelix
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/xhelix run --config /etc/xhelix/xhelix.yaml
Restart=always
RestartSec=2s
StartLimitInterval=0

# Watchdog (Phase 0: agent sends READY=1; periodic kicks land in Phase 1)
WatchdogSec=30s
NotifyAccess=main

# Capabilities. Phase 0 needs almost nothing; later phases require
# CAP_BPF, CAP_PERFMON, CAP_NET_ADMIN, CAP_SYS_PTRACE, CAP_DAC_READ_SEARCH.
# CAP_SYS_ADMIN added for fanotify FAN_OPEN_PERM (credbroker fangate, P-USG.2a).
# CAP_NET_RAW added for the dpi sensor's AF_PACKET socket — needed to
# extract TLS ClientHello SNI (P-EGRESS.M1 tuning sprint). Without it,
# SNI for outbound TLS is invisible and dest classification is IP-only.
# Honest tradeoff: CAP_NET_RAW lets a popped daemon craft raw packets;
# given the daemon already has CAP_BPF + CAP_SYS_ADMIN this is
# incremental, not categorical.
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_ADMIN CAP_NET_RAW

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/xhelix /var/log/xhelix /run/xhelix
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity
TasksMax=4096

[Install]
WantedBy=multi-user.target
